# _deploy/ircam-api.service
[Unit]
Description=IR Cam Flask API (picamera2)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dachan
WorkingDirectory=/home/dachan/ir-cam
Environment=CAMERA=real
ExecStart=/home/dachan/ir-cam/.venv/bin/python app.py --host 0.0.0.0 --port 5000
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
